# deploy/systemd/airchoice-daily-stats.timer
[Unit]
Description=AirChoice Daily Stats Timer (daily 23:10 KST)
Requires=airchoice-daily-stats.service

[Timer]
OnCalendar=*-*-* 23:10:00
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
