# deploy/systemd/airchoice-backup.timer
[Unit]
Description=AirChoice DB Backup Timer (daily 23:00 KST)
Requires=airchoice-backup.service

[Timer]
# 서버 timezone: Asia/Seoul 기준
# 확인 명령: timedatectl | grep "Time zone"
OnCalendar=*-*-* 23:00:00
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
